# etc/supervisor/conf.d/postgresql.conf
[program:postgresql]
command=/usr/bin/postgres -D /var/lib/postgresql/data

user=postgres


stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
